# bot_p2.conf
[program:bot_p2]
command=python3.10 /home/ubuntu/bot_p2/bot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/bot_p2.err.log
stdout_logfile=/var/log/bot_p2.out.log
user=root
