# bot_p.conf
[program:bot_p]
command=python3.10 /home/ubuntu/bot_p/bot.py
autostart=true
autorestart=true
stderr_logfile=/var/log/bot_p.err.log
stdout_logfile=/var/log/bot_p.out.log
user=root
